# examples/systemd/nginx-auth-subrequest-ldap.socket

[Socket]
ListenStream=/var/run/nginx-auth-subrequest-ldap.sock
User=nginx

[Install]
WantedBy=sockets.target
