# example-systemd.socket

[Socket]
ListenStream=/var/run/nginx-auth-subrequest-ldap.sock

[Install]
WantedBy=sockets.target
